# vm-udev-rules/99-nm-auto-add.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="eth*|ens*|enp*", RUN+="/usr/local/bin/nm-auto-add-connection.sh $name"